# protozoo/scripts/debian_jessie/monit/files/crontab/alive
* * * * * pzoo sh /home/spanel/protozoo/scripts/debian_jessie/monit/files/alive.sh

