# protozoo/scripts/debian_jessie/monit/files/crontab/alive
* * * * * pzoo python3 /home/spanel/protozoo/scripts/debian_jessie/monit/files/alive.py

